# dotfiles/services/user/i3wm.target
[Unit]
Description=X session managed by systemd
BindsTo=default.target
Requires=i3.service
