# dotfiles/services/user/polybar@.service
[Unit]
Description=Polybar
PartOf=default.target
After=i3.service
Requires=i3.service

[Service]
RemainAfterExit=yes
ExecStart=/usr/bin/polybar %i --config=%h/.config/polybar/config
