# dotfiles/services/user/polybar@.timer
[Unit]
Description=Timer for Polybar

[Timer]
OnStartupSec=1sec

[Install]
WantedBy=default.target
