# dotfiles/services/user/volnoti.service
[Unit]
Description=Volnoti
PartOf=default.target

[Service]
RemainAfterExit=yes
ExecStart=/usr/bin/volnoti

[Install]
WantedBy=default.target
