# dotfiles/services/light-locker.service
[Unit]
Description=Light Locker

[Service]
ExecStart=/usr/bin/light-locker &

[Install]
WantedBy=multi-user.target
