# dotfiles/services/redshift.service
[Unit]
Description=Redshift
After=graphical.target

[Service]
ExecStart=/usr/bin/redshift &

[Install]
WantedBy=multi-user.target
