# dotfiles/services/volnoti.service
[Unit]
Description=Volnoti
After=graphical.target

[Service]
ExecStart=/usr/bin/volnoti &

[Install]
WantedBy=multi-user.target
