# dotfiles/services/volnoti.service
[Unit]
Description=Volnoti

[Service]
ExecStart=/usr/bin/volnoti

[Install]
WantedBy=multi-user.target
